# tests/modules/infra/conf/manager/rancher-backup.timer
[Unit]
Description=Backup Rancher Control Plane

[Timer]
OnCalendar=*-*-* 8:00:00

[Install]
WantedBy=timers.target